# template/sample/back/npu-fw-v1.6/lib/systemd/system/rockchip-npu.service
[Unit]
Description=Rockchip NPU service
#Requires=display-manager.service
#After=display-manager.service
[Service]
Type=simple
ExecStart=/etc/init.d/npuservice start
StandardOutput=null
[Install]
WantedBy=multi-user.target
